# init/pcmkc-dped@.service
[Unit]
Description=Deployable Policy Engine
Requires=pcmkc-vmlauncher.service

[Service]
ExecStart=/usr/sbin/dped -v %i
Type=simple
