# src/dped@.service
[Unit]
Description=Deployable Policy Engine

[Service]
ExecStart=/usr/sbin/dped -v -v %i
Type=simple
